# system_files/usr/lib/systemd/system/tblue-hibernate-setup.service

[Unit]
Description=TBlue Hibernate first-boot provisioning
After=local-fs.target

[Service]
Type=oneshot
ExecStart=/usr/libexec/tblue-hibernate-setup
RemainAfterExit=yes
TimeoutSec=0
StandardOutput=journal+console
StandardError=journal

[Install]
WantedBy=multi-user.target
