# system_files/usr/lib/systemd/system/tblue-hibernate-setup.service
[Unit]
Description=TBlue Hibernate first-boot setup
After=local-fs.target
ConditionPathExists=!/var/lib/tblue/hibernate-setup.done

[Service]
Type=oneshot
ExecStart=/usr/libexec/tblue-hibernate-setup

[Install]
WantedBy=multi-user.target
